# bin/services/userclean.service
[Unit]
Description=Run Task as Shutdown
DefaultDependencies=no
Before=shutdown.target

[Service]
Type=oneshot
ExecStart=/usr/bin/clean-cacheplus.sh
TimeoutStartSec=0

[Install]
WantedBy=shutdown.target
